# arch/dot_local/chezmoi/root/etc/systemd/system/ddns-go.service.d/proxy.conf
[Service]
EnvironmentFile = /etc/environment.d/proxy.conf
