# ubuntu/dot_config/systemd/user/ddns.timer
# -*- mode: ini; -*-

[Unit]
Description = DDNS

[Install]
WantedBy = timers.target

[Timer]
OnStartupSec = 0
OnUnitActiveSec = 1min
