# modules/base-user/dot_config/systemd/user/ufw-geomap-publish.service
# -*- mode: ini; -*-

[Unit]
Description = Publish UFW GeoMap to GitHub Pages

[Service]
ExecStart = %h/.local/bin/ufw-geomap-publish.sh
