# modules/network/ddns/dot_config/systemd/user/ddns.service
# -*- mode: ini; -*-

[Unit]
Description = DDNS

[Service]
ExecStart = %h/.local/bin/ddns update --quiet
TimeoutSec = 1min
