# modules/net/mihomo/dot_config/systemd/user/mihomo-push.timer
# -*- mode: ini; -*-

[Unit]
Description = Push Mihomo Configuration

[Install]
WantedBy = timers.target

[Timer]
OnCalendar = daily
Persistent = true
