# modules/net/ddns/dot_config/systemd/user/ddns.timer
# -*- mode: ini; -*-

[Unit]
Description = DDNS

[Install]
WantedBy = timers.target

[Timer]
OnCalendar = minutely
