# modules/access/tsinghua/^etc/systemd/system/mnt-DATA41.automount
# -*- mode: ini; -*-

[Install]
WantedBy = multi-user.target

[Automount]
Where = /mnt/DATA41
