# modules/access/tsinghua/^etc/systemd/system/mnt-tsinghua.automount
# -*- mode: ini; -*-

[Install]
WantedBy = multi-user.target

[Automount]
Where = /mnt/tsinghua
TimeoutIdleSec = 10min
