# cachyos/dot_cache/exact_dotfiles/exact_root/exact_etc/exact_systemd/exact_system/exact_gdm.service.d/99-override.conf.tmpl
# -*- mode: ini; -*-

[Service]
ExecStartPre = -install --backup -D --mode="u=rw,go=r" --preserve-timestamps --no-target-directory --verbose '{{ .chezmoi.homeDir }}/.config/monitors.xml' '/var/lib/gdm/.config/monitors.xml'
