# GPUServer-402-1/dot_config/systemd/user/proxy.service
[Install]
WantedBy=default.target

[Service]
ExecStart=zsh %h/.oh-my-zsh/custom/plugins/proxy/proxy.plugin.zsh

[Unit]
Description=Auto Proxy
