# modules/storage/disk/^etc/systemd/system/mnt-pc06\x2dsda1.automount
# -*- mode: ini; -*-

[Install]
WantedBy = multi-user.target

[Automount]
Where = /mnt/pc06-sda1
TimeoutIdleSec = 10min
