# modules/storage/disk/^etc/systemd/system/mnt-PC06sda1.automount
# -*- mode: ini; -*-

[Install]
WantedBy = multi-user.target

[Automount]
Where = /mnt/PC06sda1
TimeoutIdleSec = 10min
