# modules/storage/disk/^etc/systemd/system/mnt-PC06sda1.mount.tmpl
# -*- mode: ini; -*-

[Mount]
What = /dev/disk/by-uuid/c2295d22-106c-4fbe-9588-ec3f72152381
Where = /mnt/PC06sda1
Type = xfs
Options = noatime,nodev,rw
