# arch/dot_config/systemd/user/ddns.service
# -*- mode: ini; -*-

[Unit]
Description = DDNS

[Service]
ExecStart = ddns update --quiet
TimeoutSec = 1min
