# arch/dot_config/systemd/user/restic.timer
[Install]
WantedBy = timers.target

[Timer]
OnStartupSec = 0
OnUnitActiveSec = 1h

[Unit]
Description = Restic Backup
