# arch/dot_config/systemd/user/restic.service
[Service]
ExecStart = bash %h/.local/bin/restic-backup.sh
Restart = on-failure

[Unit]
Description = Restic Backup
