# arch/dot_config/systemd/user/sing.service
# -*- mode: ini; -*-

[Unit]
Description = Update sing-box Configuration

[Service]
ExecStart = sing.sh
Restart = on-failure
RestartSec = 5s
