# arch/other/root/etc/cron.daily/executable_geosite.sh
#!/bin/bash
set -o errexit
set -o nounset
set -o pipefail

rm /var/lib/sing-box/geosite.db
systemctl restart sing-box.service
