# arch/other/root/etc/cron.daily/executable_trash-empty.sh.tmpl
#!/bin/bash
set -o errexit
set -o nounset
set -o pipefail

sudo --user="{{ .chezmoi.username }}" trash-empty 30
